# configs/gunicorn_limonkg.socket
[Unit]
Description=gunicorn socket for limonkg project

[Socket]
ListenStream=/run/gunicorn/limonkg_socket

[Install]
WantedBy=sockets.target
